# doc/service/log/run
#!/bin/sh
set -e

LOG=/var/log/lighttpd2/sv.log

test -d "$LOG" || mkdir -p -m2750 "$LOG" && chown root:adm "$LOG"
exec svlogd -tt "$LOG"
